# supervisord.conf
[supervisord]
nodaemon=true

[program:sickbeard]
command=/usr/bin/python sickbeard/SickBeard.py --config=/config/sickbeard-config.ini
redirect_stderr=true
